# supervisord.conf
; Supervisor config file for dowse.
;
; For more information on the config file, please see:
; http://supervisord.org/configuration.html
;
; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=/run/dowse/supervisord.sock
chown=dowse:dowse

[supervisord]
logfile=/var/log/dowse/supervisor/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=warn                ; log level; default info; others: debug,warn,trace
pidfile=/run/dowse/supervisord.pid
nodaemon=false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200
; supervisord has to run as root
; https://stackoverflow.com/a/36199329
environment=R="%(ENV_R)s"

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/dowse/supervisord.sock
prompt=dowse

[program:bridge_setup]
command=%(ENV_R)s/daemons/bridge_setup
priority=150
startsecs=0
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/bridge_setup.log

[program:dhcpd]
command=%(ENV_R)s/daemons/dhcpd
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/dhcpd.log

[program:dnscrypt-proxy]
command=%(ENV_R)s/daemons/dnscrypt-proxy
priority=300
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/dnscrypt-proxy.log

[program:hblock]
command=%(ENV_R)s/daemons/hblock
priority=150
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/hblock.log


[program:hostapd]
command=%(ENV_R)s/daemons/hostapd
priority=200
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/hostapd.log

[program:mosquitto]
command=%(ENV_R)s/daemons/mosquitto
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/mosquitto.log

[program:dowse-to-mqtt]
command=%(ENV_R)s/daemons/dowse-to-mqtt
; dowse-to-mqtt needs to start after mosquitto
priority=110
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/dowse-to-mqtt.log

[program:netdata]
command=%(ENV_R)s/daemons/netdata
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/netdata.log

[program:pendulum_nmap]
command=%(ENV_R)s/daemons/pendulum_nmap
priority=150
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/pendulum_nmap.log

[program:update-psl]
command=%(ENV_R)s/daemons/update-psl
priority=150
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/update-psl.log

[program:pgld]
command=%(ENV_R)s/daemons/pgld
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/pgld.log

[program:pglcmd]
command=%(ENV_R)s/daemons/pglcmd
priority=150
startsecs=0
autorestart=unexpected
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/pglcmd.log

[program:redis]
command=%(ENV_R)s/daemons/redis
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/redis.log

[program:redis-volatile]
command=%(ENV_R)s/daemons/redis-volatile
priority=110
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/redis-volatile.log

[program:unbound]
command=%(ENV_R)s/daemons/unbound
priority=100
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/unbound.log

[program:webui2]
command=%(ENV_R)s/daemons/webui2
priority=120
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/webui2.log

[program:webui2-fifo]
command=%(ENV_R)s/daemons/webui2-fifo
priority=120
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/dowse/supervisor/webui2-fifo.log

